# ansible/roles/nexus/templates/cert-renewer@.timer.j2
[Unit]
Description=Certificate renewal timer for %i

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
